# deploy/systemd/user/voc-updater.path
[Unit]
Description=Watch VOC update packages

[Path]
PathExistsGlob=/home/kasp/Project/voc_project/updates/*.tar.gz

[Install]
WantedBy=default.target
